# deploy/systemd/assessment-agent.service
[Unit]
Description=Assessment Agent — collector + worker (task.install)
# SYSTEM unit (systemctl, root). Installed under
#   /etc/systemd/system/assessment-agent.service
# and run as root so the collector can read every process' /proc/<pid>/exe,
# /proc/<pid>/fd and comm — otherwise listen_ports[].pid/comm and services[].exe
# for daemons owned by other users (sshd/nginx/mysql/…) come back null.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root

# State under /var/lib (system). StateDirectory= creates /var/lib/agent-worker
# 0700 — this matches WORKER_STATE_DIR in agent.env.example so the worker writes
# where the unit provisions. RuntimeDirectory= creates /run/assessment-agent.
StateDirectory=agent-worker
StateDirectoryMode=0700
RuntimeDirectory=assessment-agent
RuntimeDirectoryMode=0700
WorkingDirectory=/var/lib/agent-worker

# Environment — system config under /etc. Both files are optional (leading '-');
# secrets live in agent.env.local.
EnvironmentFile=-/etc/assessment-agent/agent.env
EnvironmentFile=-/etc/assessment-agent/agent.env.local

ExecStart=/usr/local/bin/assessment-agent

# C5: drain in-flight install on SIGTERM. install.timeout_sec is normally
# <= 600s, so 15 min covers the maximum drain plus result publish + ack.
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=15min

# Restart policy — survive transient broker outages and collector errors.
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

# NOTE: sandbox hardening (ProtectProc, ProtectSystem, PrivateDevices,
# ProtectKernelTunables, RestrictNamespaces, …) is deliberately NOT set. The
# agent's whole job is to read other processes' /proc entries and host /sys —
# those directives would hide exactly the data the collector must gather. The
# agent runs as root by design; containment is out of scope for this unit.
NoNewPrivileges=no

# Logging — stderr captured by the system journal (journalctl -u).
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assessment-agent

[Install]
WantedBy=multi-user.target
